# deploy/systemd/ai-image-workshop-update.path
[Unit]
Description=Watch for AI Image Workshop update requests

[Path]
PathExists=/var/lib/ai-image-workshop-updater/inbox/request.json
Unit=ai-image-workshop-update.service

[Install]
WantedBy=multi-user.target
